# deploy/thefootnotes.service
[Unit]
Description=thefootnotes — телеграм-бот и дашборд
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=thefootnotes
Group=thefootnotes
WorkingDirectory=/opt/thefootnotes
# .env читается самим приложением (load_dotenv) из WorkingDirectory.
ExecStart=/opt/thefootnotes/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

# --- Харденинг (изоляция процесса, близко к контейнерной) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Единственное место, куда сервису можно писать — папка с базой и файлами.
ReadWritePaths=/opt/thefootnotes/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
